# conf/etc/init/fix-brightness.conf
description "Sets brightness after graphics device is loaded"

start on graphics-device-added
task
exec echo 60 > /sys/class/backlight/intel_backlight/brightness
